# logmonitor.conf
description "Homebrew log monitoring daemon"
author "Georg Rempfer"

respawn
respawn limit 10 5
start on started mysql
exec /path/to/logmonitor.py
